# extras/rfz-sync.service
[Unit]
Description=rfz document synchronization service

[Service]
Type=oneshot
ExecStart=rfz sync
